# pi/systemd/kiosk-chromium.service
move to /etc/systemd/system/kiosk-chromium.service

[Unit]
Description=Launch Chromium in fullscreen kiosk mode
After=graphical.target barkasse-ui.service
Wants=barkasse-ui.service

[Service]
User=pi
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/usr/bin/chromium-browser \
  --kiosk --incognito --noerrdialogs --disable-infobars \
  --app=http://localhost:8080
Restart=always

[Install]
WantedBy=graphical.target
