# pi/systemd/kiosk-chromium.service
move to /etc/systemd/system/kiosk-chromium.service

[Unit]
Description=Launch Chromium in fullscreen kiosk mode
After=graphical.target barkasse-ui.service
Wants=barkasse-ui.service

[Service]
User=hub
Environment=XDG_RUNTIME_DIR=/run/user/1001
ExecStart=/usr/bin/chromium-browser \
  --kiosk --incognito --noerrdialogs --disable-infobars \
  --app=http://localhost:8080 \
  --start-fullscreen --check-for-update-interval=31536000
Restart=on-failure
RestartSec=5s


[Install]
WantedBy=graphical.target
